# deployment/sensor/sensor-supervisord.conf.j2
[program:otbox]
command     = /usr/bin/python /opt/opentestbed/otbox.py -v --testbed "wilab" --broker "{{ broker_hostname }}"
autostart   = true
autorestart = true
directory = /opt/opentestbed