# deployment/sensor/sensor-supervisord.conf.j2
[program:otbox]
command     = /usr/local/bin/opentestbed -v
environment = OTB_TESTBED="wilab", OTB_BROKER="{{ broker_hostname }}"
autostart   = true
autorestart = true
directory = /tmp